# cron/launchd/org.oracle.memory-archive-timeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.oracle.memory-archive-timeline</string>
    <!--
      Monthly Events Timeline rollover. On day 1 of each month, moves any
      MEMORY.md ### YYYY-MM section earlier than the current month into
      memory/timeline-archive.md (handles overdue catch-up automatically
      — not just "last month"). Pure shell job; no LLM, no runner.sh.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd __PROJECT_DIR__ &amp;&amp; bash scripts/cron-memory-archive.sh timeline &gt;&gt; __PROJECT_DIR__/cron/logs/memory-archive-launchd.log 2&gt;&amp;1</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>10</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/cron/logs/memory-archive-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/cron/logs/memory-archive-launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>
