# cron/launchd/org.oracle.memory-janitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.oracle.memory-janitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/cron/runner.sh</string>
        <string>memory-janitor</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>7</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/cron/logs/memory-janitor-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/cron/logs/memory-janitor-launchd.log</string>
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>
